# services/syncthing.service
[Unit]
Description=Syncthing service

[Service]
ExecStart=/usr/bin/syncthing serve

[Install]
WantedBy=default.target
